# onlinejp.conf
[program:onlinejp]
command = /srv/onlinejp/gunicorn_start.bash
user = root
stdout_logfile = /srv/onlinejp/log/gunicorn_supervisor.log
redirect_stderr = true
